# cookbooks/zookeeper/templates/default/upstart.conf.erb
# Start the <%= @service %> server daemon
#
author          "Librato <dev@librato.com>"
description     "Runs the <%= @service %> server daemon"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

env JAVA_HOME="<%= @java_home %>"
env CLASSPATH="<%= @node[:zookeeper][:conf_dir] %>:/opt/zookeeper/*:/opt/zookeeper/lib/*"
env ZOOMAIN="org.apache.zookeeper.server.quorum.QuorumPeerMain"
env ZOOCFG="<%= @node[:zookeeper][:conf_dir] %>/zoo.cfg"

env JAVA_OPTS="-Xms<%= @node[:zookeeper][:min_heap] %> -Xmx<%= @node[:zookeeper][:max_heap] %> -Dcom.sun.management.jmxremote.port=<%= @node[:zookeeper][:jmx_port] %> -Dcom.sun.management.jmxremote.authenticate=false -Dcom.sun.management.jmxremote.ssl=false"

setuid zookeeper

script
        exec $JAVA_HOME/bin/java $JAVA_OPTS -cp $CLASSPATH $JVMFLAGS $ZOOMAIN $ZOOCFG
end script
